# supervisord.conf
[program:cool]
command = /home/dotcloud/current/script/dotcloud.sh
